# harvester/root/etc/services.d/chiaharvester/finish
#!/usr/bin/with-contenv sh
export HOME=/chia
rm ${HOME}/.chia/mainnet/run/chia_harvester.pid

