# harvester/root/etc/services.d/chiadaemon/run
#!/usr/bin/with-contenv sh
export HOME=/chia
exec s6-setuidgid chia chia run_daemon

